# glpi-official/services/web/files/etc/cron.d/glpi
* * * * *    /usr/local/bin/php /var/www/html/glpi/front/cron.php &>/dev/null
